# docker/supervisord-s3-ftp.conf
[program:s3-ftp]
command=/start-s3-ftp.sh
numprocs=1
autostart=true
autorestart=true